# pkg/dpkg/postinst
#!/bin/bash

udevadm control --reload-rules

systemctl daemon-reload
systemctl enable midisense
systemctl start midisense
