# pkg/dpkg/prerm
#!/bin/bash

systemctl stop midisense
systemctl disable midisense
