# src/34-usb-midisense.rules
ACTION=="add|remove", SUBSYSTEM=="usb", RUN+="/usr/bin/midisense.py -u"

