# src/45-bt-midisense.rules
ACTION=="add|remove", SUBSYSTEM=="bluetooth", RUN+="/usr/bin/midisense.py -u"

